# forms/dockerconfigs/conf.d/python.conf
[program:python]
command="python" "/app/main.py"
autorestart=true
