# forms/dockerconfigs/conf.d/xterm.conf
[program:xterm]
command=_xterm
autorestart=false
